# borg-backup.service
[Service]
type=oneshot
ExecStart=/etc/borg/borg-run.sh
